# navimow-observability.service
[Unit]
Description=Navimow Observability Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/data/navimow-observability/backend
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5
EnvironmentFile=/data/navimow-observability/backend/.env

[Install]
WantedBy=multi-user.target
